# server_config/systemd/sidekiq-long@.service
# Based on https://github.com/mperham/sidekiq/blob/master/examples/systemd/sidekiq.service
#
# Templated unit: run N independent long_update workers, one OS process per
# instance. Enable, e.g., three of them with:
#
#   sudo systemctl enable --now sidekiq-long@1 sidekiq-long@2 sidekiq-long@3
#
# Each instance is a separate process with its own GVL, so they run on separate
# cores. A single `--concurrency 3` process cannot: MRI serializes Ruby on one
# GVL per process, so it pegs ~1 core (leaving the rest idle) and starves its
# own heartbeat thread — which is what produces the bogus "Redis network
# connection appears to be performing poorly" warnings.

[Unit]
Description=sidekiq-long@%i
After=syslog.target network.target

[Service]
Type=simple
WorkingDirectory=/home/ragesoss/WikiEduDashboard
ExecStart=/home/ragesoss/.rvm/bin/rvm default do bundle exec sidekiq -e production --queue long_update --queue short_update --queue medium_update --concurrency 1
User=root
Group=root
UMask=0002

# Greatly reduce Ruby memory fragmentation and heap usage
# https://www.mikeperham.com/2018/04/25/taming-rails-memory-bloat/
Environment=MALLOC_ARENA_MAX=2

# if we crash, restart
RestartSec=1
Restart=on-failure

# output goes to /var/log/syslog
StandardOutput=syslog
StandardError=syslog

SyslogIdentifier=sidekiq-long-%i

[Install]
WantedBy=multi-user.target
